# crontab/restart-kiosk
*/5 * * * * ### /home/###/restart-kiosk.sh &
